# platform/broadcom/sonic-platform-modules-micas/m2-w6940-128x1-fr4-pure-pd/cpo/service/cpo_daemon.service
[Unit]
Description=cpo_daemon service
After=platform_driver.service

[Service]
ExecStart=/usr/local/bin/cpo_daemon.sh start
ExecStop=/usr/local/bin/cpo_daemon.sh stop
KillMode=process
Restart=on-failure
RestartSec=60s

[Install]
WantedBy=multi-user.target
