# platform/broadcom/sonic-platform-modules-micas/m2-w6940-128x1-fr4-pure-pd/cpo/service/oe_voltage_setting.service
[Unit]
Description=
Requires=platform_driver.service
After=platform_driver.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/oe_voltage_setting.py

[Install]
WantedBy=multi-user.target