# files/image_config/systemd/set-epoch.timer
[Unit]
Description=Update the epoch file with the current time daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
